# duck-control.service
[Unit]
Description=ChatGPT Duck Control Panel
After=network.target

[Service]
Type=simple
User=admog
WorkingDirectory=/home/admog/Code/chatgpt-and-pi5/chatgpt-and
ExecStart=/home/admog/Code/chatgpt-and-pi5/chatgpt-and/.venv/bin/python /home/admog/Code/chatgpt-and-pi5/chatgpt-and/duck-control.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
