# services/fan-control.service
[Unit]
Description=ChatGPT Duck Fan Control
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/admog/Code/chatgpt-and
ExecStart=/home/admog/Code/chatgpt-and/.venv/bin/python /home/admog/Code/chatgpt-and/fan_control.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
